# dashboard.service
[Unit]
Description=Dashboard (FastAPI)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/www

ExecStart=/root/.local/bin/uv run python -m backend --host 0.0.0.0 --port 8000 --log-level info

Restart=always
RestartSec=3

StandardOutput=journal
StandardError=journal

Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target